# pi-cec-daemon.service
[Unit]
Description=Pi CEC Daemon - HDMI CEC automation for TV and peripherals
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/pi-cec-daemon
ExecStart=/opt/pi-cec-daemon/venv/bin/python3 /opt/pi-cec-daemon/cec_daemon.py /opt/pi-cec-daemon/config.yaml
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pi-cec-daemon

# Security hardening (optional but recommended)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
